# scripts/zamoritos.service
[Unit]
Description=Zamoritos – Stack Docker Compose
Documentation=https://docs.docker.com/compose/
After=network-online.target docker.socket
Wants=network-online.target docker.socket

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/juan/zamoritos
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
