# scripts/zamoritos.service
[Unit]
Description=Zamoritos – Stack Docker Compose
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=juan
WorkingDirectory=/home/juan/zamoritos
ExecStart=/home/juan/zamoritos/scripts/zamoritos-start.sh
ExecStop=/bin/bash -c '\
    export PATH="$HOME/.docker/bin:/usr/local/bin:/usr/bin:/bin"; \
    cd /home/juan/zamoritos && \
    docker compose down || docker-compose down'
TimeoutStartSec=120
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
