# scripts/zamoritos-backup.service
[Unit]
Description=Zamoritos – Backup de base de datos (mysqldump)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/juan/zamoritos
ExecStartPre=/bin/mkdir -p /home/juan/zamoritos/backend/storage/backups
ExecStart=/bin/bash -c '\
    FILE=/home/juan/zamoritos/backend/storage/backups/db_$(date +%%Y%%m%%d_%%H%%M%%S).sql && \
    /usr/bin/docker compose exec -T db \
        mysqldump -u zamoritos -pzamoritos_secret zamoritos > "$FILE" && \
    gzip "$FILE" && \
    echo "Backup OK: ${FILE}.gz"'
ExecStartPost=/bin/find /home/juan/zamoritos/backend/storage/backups \
    -name "*.sql.gz" -mtime +30 -delete
StandardOutput=journal
StandardError=journal
